# docker/php/crontab
* */3 * * * cd /var/www/application && php artisan challenges:reset >> /var/log/cron.log 2>&1
